# templates/default/consul.conf.erb
<% @agent   = node.consul.agent %>
<% @upstart = @agent.upstart %>

author "Brandon Dennis <me@brandondennis.me>"

env GOMAXPROCS=<%= @upstart.gomaxprocs || 1 %>
export GOMAXPROCS

start on runlevel [2345]
stop on runlevel [!2345]

<% if @upstart.respawn.enable %>
  <% @respawn = @upstart.respawn %>
respawn
respawn limit <%= "#{@respawn.limit} #{@respawn.duration}" %>
<% end %>

reload signal SIGHUP

exec consul agent -config-dir <%= @agent.config.config_dir %>
